# provisioning/roles/firewall/templates/firewall.upstart.conf.j2
description "Start firewall"

start on net-device-up IFACE=lo
console output

pre-start exec /etc/init.d/local/firewall start
post-start exec /etc/init.d/local/firewall save
